# .config/autostart/arbtt-capture.desktop

[Desktop Entry]
Type=Application
Name=arbtt
Exec=arbtt-capture
Icon=system-run
Comment=
