# .config/autostart/davmail.desktop

[Desktop Entry]
Type=Application
Name=Davmail
Exec=davmail
Icon=system-run
Comment=Local portal to exchange server
