# .config/autostart/firefox-3.5.desktop

[Desktop Entry]
Type=Application
Name=Firefox
Exec=firefox-3.5
Icon=system-run
Comment=
