# .config/autostart/gvim.desktop

[Desktop Entry]
Type=Application
Name=Gvim
Exec=gvim
Icon=system-run
Comment=
