# .config/autostart/sunbird.desktop

[Desktop Entry]
Type=Application
Name=Sunbird
Exec=sunbird
Icon=system-run
Comment=Calendering
